# otel/agent/otelcol-agent.container
[Unit]
Description=OpenTelemetry Collector agent (CNCF Incubating)
After=network-online.target

[Container]
Image=docker.io/otel/opentelemetry-collector-contrib:0.111.0
ContainerName=otelcol-agent
# Config
Volume=/etc/otelcol/config.yaml:/etc/otelcol-contrib/config.yaml:ro,z
# Systemd journal logs
Volume=/run/log/journal:/run/log/journal:ro
# Host metrics (CPU, memory, disk, network)
Volume=/proc:/hostfs/proc:ro
Volume=/sys:/hostfs/sys:ro
Volume=/:/hostfs:ro
SecurityLabelDisable=true
PodmanArgs=--pid=host

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
